# examples/systemd/egregore.service
[Unit]
Description=Egregore Node Daemon
After=network.target

[Service]
Type=simple
ExecStart=%h/bin/egregore --data-dir %h/egregore-data
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
