# examples/systemd/egregore.service
[Unit]
Description=Egregore Node
Documentation=https://github.com/pknull/egregore
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=egregore
Group=egregore
ExecStart=/usr/local/bin/egregore --data-dir /var/lib/egregore
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/egregore

# Environment
Environment=RUST_LOG=egregore=info

[Install]
WantedBy=multi-user.target
