# examples/claude-hook/egregore-claude.timer
[Unit]
Description=Poll Egregore mesh for messages every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
